# skills/code-review-toolkit/code-review-mcp.service
[Unit]
Description=Code Review MCP Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/.openclaw/workspace/team3/code-review-mcp
Environment=PATH=/home/ubuntu/.local/bin:/home/ubuntu/.nvm/versions/node/v22.23.0/bin:/usr/local/bin:/usr/bin:/bin:/home/ubuntu/.cargo/bin:/home/ubuntu/.foundry/bin
ExecStart=/usr/bin/python3 server.py 9020
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
